# time_sync.service
[Unit]
Description=Time Sync Service
After=network.target

[Service]
ExecStart=/usr/bin/time_sync.sh

[Install]
WantedBy=default.target
